# config/systemd/cows.all.target
[Unit]
Description=Start/stop all services
After=default.target

[Install]
RequiredBy=default.target